# user/packages/src-backup/dot-local/share/systemd/user/src-backup.service
[Unit]
Description=snapshot src/

[Service]
ExecStart=%h/.local/libexec/src-backup
Restart=always
RestartSec=1h

[Install]
WantedBy=default.target
